# xremap/.config/systemd/user/xremap.service
[Unit]
Description=Xremap Daemon

[Service]
ExecStart=/usr/bin/xremap /home/yayu/.config/xremap/config.yml

[Install]
WantedBy=default.target